# templates/ryu.conf.erb
description "Ryu Server upstart job"
author "Andy Edmonds <andy@edmonds.be>"

start on runlevel [2345]
stop on runlevel [!2345]


chdir <%= scope.lookupvar("ryu::params::sys_rundir") %>

pre-start script
    mkdir -p <%= scope.lookupvar("ryu::params::ryu_rundir") %>
    chown ryu:root <%= scope.lookupvar("ryu::params::ryu_rundir") %>

    # mkdir -p <%= scope.lookupvar("ryu::params::ryu_lockdir") %>
    # chown ryu:root <%= scope.lookupvar("ryu::params::ryu_lockdir") %>
end script

exec start-stop-daemon --start --chuid ryu --exec <%= scope.lookupvar("ryu::params::ryu_service") %> -- --config-file=<%= scope.lookupvar("ryu::params::ryu_conf_file") %>